# bin/ethereum_management/geth.service
[Unit]
Description=Academia

[Service]
Type=simple
User=root
Restart=always
WorkingDirectory=/root/academia_network
ExecStart=/usr/bin/geth --config academianodeconf.toml

[Install]
WantedBy=default.target